# staging/src/k8s.io/apiserver-builder/vendor/k8s.io/kubernetes/cluster/saltbase/salt/supervisor/docker.conf
[program:docker]
command=/usr/sbin/docker-checker.sh
stderr_logfile=/var/log/supervisor/docker-stderr.log
stdout_logfile=/var/log/supervisor/docker-stdout.log
autorestart=true
startretries=1000000
