# staging/src/k8s.io/apiserver-builder/vendor/k8s.io/kubernetes/cluster/saltbase/salt/supervisor/kube-addons.conf
[program:kube-addons]
command=/usr/sbin/kube-addons-checker.sh
stderr_logfile=/var/log/supervisor/kube-addons-stderr.log
stdout_logfile=/var/log/supervisor/kube-addons-stdout.log
autorestart=true
startretries=1000000
