# vendor/github.com/coreos/rkt/stage1/units/units/exit.target
[Unit]
Description=Exit the container
DefaultDependencies=no
AllowIsolate=yes
