# vendor/github.com/coreos/rkt/stage1/units/units/halt.target
[Unit]
Description=Halt
DefaultDependencies=no
AllowIsolate=yes
