# vendor/github.com/coreos/rkt/stage1/units/units/local-fs.target
[Unit]
Description=Hook into early systemd for socket-activated systemd instances
DefaultDependencies=false
Requires=sockets.target
