# vendor/github.com/coreos/rkt/stage1/units/units/poweroff.target
[Unit]
Description=Poweroff
DefaultDependencies=no
AllowIsolate=yes
