# vendor/github.com/coreos/rkt/stage1/units/units/reboot.target
[Unit]
Description=Reboot
DefaultDependencies=no
Requires=systemd-reboot.service
After=systemd-reboot.service
AllowIsolate=yes
