# vendor/github.com/coreos/rkt/stage1/units/units/sockets.target
[Unit]
Description=Sockets
DefaultDependencies=false
