# vendor/github.com/coreos/rkt/stage1/units/units/systemd-reboot.service
[Unit]
Description=Reboot
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --force reboot
