# vendor/github.com/coreos/rkt/stage1/units/units/sysusers.service
[Unit]
Description=Create /etc/passwd and /etc/group
DefaultDependencies=false
OnFailureJobMode=fail

[Service]
Type=oneshot
Restart=no
Environment="SYSTEMD_LOG_LEVEL=warning"
ExecStart=/usr/bin/systemd-sysusers
User=0
Group=0
