# vendor/github.com/coreos/rkt/stage1/units/units/ttymux@.service
[Unit]
Description=TTY handler for %i
DefaultDependencies=false
OnFailureJobMode=fail
BindsTo=%i.service

[Service]
Type=notify
Restart=no
EnvironmentFile=/rkt/iottymux/%i/env
ExecStart=/iottymux --action=ttymux --app=%i
User=0
Group=0
SyslogIdentifier=%i
StandardOutput=journal+console
StandardError=journal+console
