# vendor/github.com/coreos/etcd/contrib/systemd/etcd2-backup-coreos/etcd2-backup.timer
[Unit]
Description=etcd2-backup service timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target
