# vendor/github.com/coreos/rkt/stage1/usr_from_coreos/manifest-arm64-usr.d/journald.manifest
lib/systemd/system/syslog.socket
lib/systemd/system/systemd-journald-audit.socket
lib/systemd/system/systemd-journald-dev-log.socket
lib/systemd/system/systemd-journald.service
lib/systemd/system/systemd-journald.socket
