# vendor/github.com/coreos/rkt/stage1/usr_from_coreos/manifest-amd64-usr.d/journald.manifest
lib64/systemd/system/syslog.socket
lib64/systemd/system/systemd-journald-audit.socket
lib64/systemd/system/systemd-journald-dev-log.socket
lib64/systemd/system/systemd-journald.service
lib64/systemd/system/systemd-journald.socket
