# vendor/github.com/heketi/heketi/extras/docker/gluster/gluster-setup.service
[Unit]
Description=Configuring GlusterFS in container
Before=rpcbind.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/gluster-setup.sh

[Install]
WantedBy=multi-user.target

