# vendor/github.com/coreos/go-systemd/fixtures/enable-disable.service
[Unit]
Description=enable disable test

[Service]
ExecStart=/bin/sleep 400
