# vendor/github.com/coreos/go-systemd/fixtures/mask-unmask.service
[Unit]
Description=mask unmask test

[Service]
ExecStart=/bin/sleep 400
