# vendor/github.com/coreos/go-systemd/fixtures/start-stop.service
[Unit]
Description=start stop test

[Service]
ExecStart=/bin/sleep 400
