# staging/src/k8s.io/apiserver-builder/vendor/k8s.io/kubernetes/cluster/ubuntu/minion/init_conf/kube-proxy.conf
description "Kube-Proxy service"
author "@jainvipin"

respawn

# start in conjunction with flanneld
start on started flanneld
stop on stopping flanneld

limit nofile 65536 65536

pre-start script
	# see also https://github.com/jainvipin/kubernetes-start
	KUBE_PROXY=/opt/bin/$UPSTART_JOB
	if [ -f /etc/default/$UPSTART_JOB ]; then
		. /etc/default/$UPSTART_JOB
	fi
	if [ -f $KUBE_PROXY ]; then
		exit 0
	fi
    exit 22
end script

script
	# modify these in /etc/default/$UPSTART_JOB (/etc/default/kube-proxy)
	KUBE_PROXY=/opt/bin/$UPSTART_JOB
	KUBE_PROXY_OPTS=""
	if [ -f /etc/default/$UPSTART_JOB ]; then
		. /etc/default/$UPSTART_JOB
	fi
	exec "$KUBE_PROXY" $KUBE_PROXY_OPTS
end script
