# cluster/saltbase/salt/supervisor/kubelet.conf
[program:kubelet]
command=/usr/sbin/kubelet-checker.sh
stderr_logfile=/var/log/supervisor/kubelet-stderr.log
stdout_logfile=/var/log/supervisor/kubelet-stdout.log
autorestart=true
startretries=1000000
