# staging/src/k8s.io/apiserver-builder/vendor/k8s.io/kubernetes/cluster/saltbase/salt/docker/docker-healthcheck.service
[Unit]
Description=Run docker-healthcheck once

[Service]
Type=oneshot
ExecStart=/opt/kubernetes/helpers/docker-healthcheck

[Install]
WantedBy=multi-user.target
