# staging/src/k8s.io/apiserver-builder/vendor/k8s.io/kubernetes/cluster/saltbase/salt/kube-master-addons/kube-master-addons.service
[Unit]
Description=Kubernetes-Master Addon Object Manager
Documentation=https://github.com/kubernetes/kubernetes

[Service]
ExecStart=/etc/kubernetes/kube-master-addons.sh

[Install]
WantedBy=multi-user.target
