# vendor/github.com/coreos/go-systemd/examples/activation/httpserver/hello.socket
[Socket]
ListenStream=127.0.0.1:8076

[Install]
WantedBy=sockets.target
