# vendor/github.com/coreos/go-systemd/examples/activation/httpserver/hello.service
[Unit]
Description=Hello World HTTP
Requires=network.target
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/httpserver

[Install]
WantedBy=multi-user.target
