# vendor/github.com/coreos/rkt/stage1/usr_from_kvm/files/sshd@.service
[Unit]
Description=OpenSSH per-connection server daemon for %I
DefaultDependencies=no
Requires=sshd-prep.service
After=sshd-prep.service

[Service]
ExecStart=-/usr/sbin/sshd -i -e
StandardInput=socket
StandardError=syslog
