# vendor/github.com/coreos/rkt/stage1/usr_from_kvm/files/sshd.socket
[Unit]
Description=OpenSSH Server Socket for entering purpose
DefaultDependencies=no

[Socket]
ListenStream=122
Accept=yes

[Install]
WantedBy=default.target
