# vendor/k8s.io/heapster/kafka/supervisor/zookeeper.conf
[program:zookeeper]
command=/usr/share/zookeeper/bin/zkServer.sh start-foreground
autostart=true
autorestart=true
