# staging/src/k8s.io/apiserver-builder/vendor/github.com/coreos/go-systemd/fixtures/reload.service
[Unit]
Description=reload test

[Service]
ExecStart=/bin/bash -c "trap '' HUP; /bin/sleep 400"
ExecReload=-/bin/systemctl kill -s HUP reload.service
