# staging/src/k8s.io/apiserver-builder/vendor/github.com/coreos/go-systemd/fixtures/start-failed.service
[Unit]
Description=starting a failed test

[Service]
ExecStartPre=/bin/false
ExecStart=/bin/sleep 400
