# vendor/github.com/coreos/rkt/dist/init/systemd/rkt-api.service
[Unit]
Description=rkt api service
Documentation=http://github.com/coreos/rkt
After=network.target rkt-api-tcp.socket
Requires=rkt-api-tcp.socket

[Service]
ExecStart=/usr/bin/rkt api-service

[Install]
WantedBy=multi-user.target
