# vendor/github.com/coreos/rkt/dist/init/systemd/rkt-gc.service
[Unit]
Description=Garbage Collection for rkt

[Service]
Environment=GRACE_PERIOD=24h
Type=oneshot
ExecStart=/usr/bin/rkt gc --grace-period=${GRACE_PERIOD}
