# vendor/github.com/coreos/rkt/dist/init/systemd/rkt-gc.timer
[Unit]
Description=Periodic Garbage Collection for rkt

[Timer]
OnActiveSec=0s
OnUnitActiveSec=12h

[Install]
WantedBy=default.target
