# vendor/github.com/coreos/rkt/dist/init/systemd/rkt-metadata.socket
[Unit]
Description=rkt metadata service socket
PartOf=rkt-metadata.service

[Socket]
ListenStream=/run/rkt/metadata-svc.sock
SocketMode=0660
SocketUser=root
SocketGroup=root
RemoveOnStop=true

[Install]
WantedBy=sockets.target
